# digests_project/systemd/kb-l2.timer

[Unit]
Description=KB L2 Builder Timer (hourly)

[Timer]
OnCalendar=hourly
Persistent=true

[Install]
WantedBy=timers.target
